# deployment/systemd_services/ai-hub-runner.service
[Unit]
Description=AI Hub Runner (APScheduler)
After=network-online.target
Requires=ai-hub-mcp.service

[Service]
Type=simple
User=reekroo
EnvironmentFile=/etc/default/ai-hub
WorkingDirectory=/home/reekroo/ai_hub
ExecStart=/home/reekroo/ai_hub/.venv_ai_hub/bin/python -m runner

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target