# deployment/systemd_services/ai-hub-mcp.socket
[Unit]
Description=AI Hub MCP Socket

[Socket]
ListenStream=/run/ai-hub-mcp.sock
SocketUser=reekroo
SocketGroup=reekroo
SocketMode=0660
Accept=yes

[Install]
WantedBy=sockets.target
