# deployment/systemd_services/ai-hub-mcp@.service
[Unit]
Description=AI Hub MCP (per-connection) Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

WorkingDirectory=/home/reekroo/ai_hub
EnvironmentFile=/etc/default/ai-hub
ExecStart=/home/reekroo/ai_hub/.venv_ai_hub/bin/python -m mcp.server

StandardInput=socket
StandardOutput=socket
StandardError=journal
Restart=on-failure

User=root
Group=root